# com.roman.apple-notes-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.roman.apple-notes-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/roman/dev/apple-notes-backup/scripts/backup-apple-notes.sh</string>
        <string>--if-stale</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>86400</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/roman/Library/Logs/apple-notes-backup-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/roman/Library/Logs/apple-notes-backup-launchd.log</string>
</dict>
</plist>
