# io.github.apple-notes-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.github.apple-notes-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/apple-notes-backup/scripts/backup-apple-notes.sh</string>
        <string>--if-stale</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>BACKUP_DIR</key>
        <string>/path/to/Backups/Apple Notes</string>
        <key>BACKUP_INTERVAL_DAYS</key>
        <string>10</string>
        <key>AUTO_CLEANUP</key>
        <string>false</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>86400</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/path/to/Library/Logs/apple-notes-backup-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/Library/Logs/apple-notes-backup-launchd.log</string>
</dict>
</plist>
